# archlinux/etc/profile.d/dircolors.sh
eval `dircolors /etc/profile.d/dircolors.256dark`
